# ansible/templates/translator_server.service.j2
[Unit]
Description=Translator Server
After=network.target

[Service]
ExecStart=/opt/translator/translator
Restart=always

[Install]
WantedBy=multi-user.target

